# deploy/systemd/agent-learn-refresh.timer
[Unit]
Description=Run Agent Learn token usage data refresh hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=5m
Unit=agent-learn-refresh.service

[Install]
WantedBy=timers.target
