# deploy/systemd/agent-learn-lark-candidates.service
[Unit]
Description=Collect daily Agent candidates from Lark messages
After=network-online.target
Wants=network-online.target
ConditionPathExists=/home/ubuntu/tool/agent_learn/package.json

[Service]
Type=oneshot
User=ubuntu
Environment=NEXT_TELEMETRY_DISABLED=1
Environment=AGENT_LEARN_SOURCE_DIR=/home/ubuntu/tool/agent_learn
Environment=AGENT_LEARN_STATE_DIR=/var/lib/agent-learn
ExecStart=/usr/bin/env bash -lc 'cd "$AGENT_LEARN_SOURCE_DIR" && npm run collect:lark-agent-candidates'
