# deploy/systemd/agent-learn-refresh.path
[Unit]
Description=Watch for Agent Learn manual refresh requests

[Path]
PathChanged=/var/lib/agent-learn/refresh.request
Unit=agent-learn-refresh.service

[Install]
WantedBy=multi-user.target
